# .conf/supervisor.ini
[unix_http_server]
file=/var/www/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/www/logs/supervisor.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/www/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[supervisorctl]
serverurl=unix:///var/www/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:demo]
user = root
directory = /var/www/demo
command = /var/www/demo/venv/bin/uwsgi /var/www/demo/.conf/uwsgi.ini
autostart = true
autorestart = true
stderr_logfile = /var/www/logs/uwsgi-err.log
stdout_logfile = /var/www/logs/uwsgi-out.log
